# systemd/homa.service
[Unit]
Description=Homa orchestrator
Documentation=https://github.com/OxQuasar/homa
After=network-online.target
Wants=network-online.target

# Burst-limit knobs live in [Unit], NOT [Service] (systemd quirk).
# Window of 300s with 15 restart attempts — generous enough that a
# brief flap (e.g. network transient at boot, podman daemon hiccup)
# doesn't cause systemd to give up and leave the daemon down.
StartLimitIntervalSec=300
StartLimitBurst=15

[Service]
# Foreground binary, no daemonization.
Type=simple
WorkingDirectory=%h/homa
ExecStart=%h/homa/homa -config config.json

# %h expands to the user's home — keeps this unit file user-relocatable
# (works for any user that has homa installed under ~/homa).

# Restart policy for the always-on orchestrator daemon.
#
#   Restart=always   auto-restart on any exit (crash, signal, even
#                    unexpected exit 0). Operator `systemctl stop`
#                    still bypasses — systemd treats explicit stops
#                    differently from process death.
#   RestartSec=5s    wait 5s between attempts (avoids tight loop
#                    during e.g. config-load failures).
#   StartLimit*      see [Unit] — burst budget for restart attempts.
Restart=always
RestartSec=5s

# Bound the graceful-shutdown window. The orchestrator's shutdownGrace
# constant is 10s; give a buffer so an in-flight provision can finish
# before SIGKILL.
TimeoutStopSec=20s

# Logs land in journalctl --user -u homa.
StandardOutput=journal
StandardError=journal

# Resource ceilings — orchestrator itself is tiny; sandboxes live in
# their own podman cgroups so these only bound the Go process.
LimitNOFILE=65536

[Install]
WantedBy=default.target
